# deploy/webhook.service
[Unit]
Description=Truckline Docs GitHub Webhook Listener
After=network.target

[Service]
Type=simple
User=truckline
Group=truckline
WorkingDirectory=/home/truckline/web_servers/docs/deploy
Environment="WEBHOOK_SECRET=CHANGE_ME_TO_A_RANDOM_SECRET"
Environment="DEPLOY_SCRIPT=/home/truckline/web_servers/docs/deploy/deploy.sh"
Environment="PORT=9000"
ExecStart=/usr/bin/python3 /home/truckline/web_servers/docs/deploy/webhook-server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
